# modules/base/files/etc/cron.d/puppet.cron
* * * * * root puppet agent --test --server puppetmaster.scup.com
